# hera/tools/shellinabox_flask/local_auto_clean/crontab
0 0 * * 0 /local_docker_container_image_clean.sh
